# etc/init.d/aweful-dns
#!/sbin/openrc-run

name="busybox $RC_SVCNAME"
command="/root/aweful-dns"
pidfile="/run/$RC_SVCNAME.pid"
command_args=""
command_background=true
output_log="/var/log/aweful-dns.log"
error_log="/var/log/aweful-dns.log"

depend() {
        need localmount
        need logger
}

start_post() {
        # restart dnsmasq service since we modified the config
        # reload doesn't work if the config was previously loaded, but was modified
        /etc/init.d/dnsmasq restart
}
